# container/interactive/IJulia/supervisord.conf
[supervisord]
nodaemon=true
logfile_maxbytes = 1MB
logfile_backups = 2

[program:shellinabox]
command=shellinaboxd -t -s /:juser:juser:/home/juser:/bin/bash --user-css Light:-/usr/local/share/doc/shellinabox/black-on-white.css,Dark:+/usr/local/share/doc/shellinabox/white-on-black.css
stdout_logfile = /home/juser/.juliabox/shellinabox.log
stdout_logfile_backups = 2
stdout_logfile_maxbytes = 1MB
stderr_logfile = /home/juser/.juliabox/shellinabox_err.log
stderr_logfile_backups = 2
stderr_logfile_maxbytes = 1MB
environment = CMDSTAN_HOME="/usr/share/cmdstan"

[program:ijulia]
command=ipython notebook
stdout_logfile = /home/juser/.juliabox/ijulia.log
stdout_logfile_backups = 2
stdout_logfile_maxbytes = 1MB
stderr_logfile = /home/juser/.juliabox/ijulia_err.log
stderr_logfile_backups = 2
stderr_logfile_maxbytes = 1MB
environment = CMDSTAN_HOME="/usr/share/cmdstan"

[program:tornado]
command=/home/juser/.juliabox/tornado/src/fmanage.py
directory=/home/juser/.juliabox/tornado
stdout_logfile = /home/juser/.juliabox/tornado/tornado.log
stdout_logfile_backups = 2
stdout_logfile_maxbytes = 1MB
stderr_logfile = /home/juser/.juliabox/tornado/tornado_err.log
stderr_logfile_backups = 2
stderr_logfile_maxbytes = 1MB
